# initramfs/buildroot/board/sunplus/common/rootfs_overlay/busybox/etc/profile.d/setLibraryPath.sh
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/usr/lib/vicore_ndk/
export GST_INSPECT_NO_COLORS=1
export PS1='\u@\h:\w\$'
